# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=python app.py
autostart=true
autorestart=true

[program:bot]
command=python bot.py
autostart=true
autorestart=true
